# harbour-countdown.timer
[Unit]
Description=Countdown timer

[Timer]
OnBootSec=200
OnUnitInactiveSec=3600

[Install]
WantedBy=timers.target

