# harbour-countdown.timer
[Unit]
Description=Countdown timer

[Timer]
OnBootSec=200
OnUnitActiveSec=120

[Install]
WantedBy=multi-user.target

