# harbour-countdown.timer
[Unit]
Description=Countdown timer

[Timer]
OnBootSec=200
OnUnitInactiveSec=60

[Install]
WantedBy=timers.target

